# server/systemd/assign-client-names.timer
[Unit]
Description=Run assign-client-names daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
