# ops/voice-ha/superparty-voice.service
[Unit]
Description=Superparty Voice Service (Docker Compose)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/opt/superparty/ops/voice-ha
ExecStart=/usr/bin/docker compose -f docker-compose.voice.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.voice.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
