# infra/orchestrator/voip_orch.service
[Unit]
Description=VoIP Remediation Orchestrator
After=network.target

[Service]
User=deploy
Environment=ORCH_SSH_USER=superparty
Environment=ORCH_SSH_HOST=46.225.182.127
Environment=DOMAIN=https://api.superparty.ro
Environment=ADMIN_TOKEN=<REDACTED_FROM_REPO>
Environment=SLACK_WEBHOOK_URL=<REDACTED_FROM_REPO>
Environment=ORCH_SECRET=<REDACTED_FROM_REPO>
WorkingDirectory=/opt/voip-orch
ExecStart=/usr/bin/env uvicorn orchestrator:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=10
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
